# crontab
0 0 * * * python /code/manage.py update_maturi_status >> /var/log/cron.log 2>&1
